# artifacts/overrides_nix/etc/profile.d/nix-daemon.sh
# Add Nix to the system profile
# this will put nix in the $PATH
. /nix/var/nix/profiles/default/etc/profile.d/nix-daemon.sh
