# artifacts/overrides_kuberblue/etc/systemd/system/kubelet.service.d/kuberblue.conf
[Service]
# On first boot, kubelet starts before kubeadm init has created
# /var/lib/kubelet/config.yaml. Without this override, kubelet hits
# systemd's default StartLimitBurst (5 starts / 10s) and enters hard-failed
# state, which prevents kubeadm init from being able to start it.
# Setting StartLimitIntervalSec=0 disables the burst limit entirely.
# RestartSec=10 slows the retry loop so it doesn't spin-burn CPU.
StartLimitIntervalSec=0
RestartSec=10
