# artifacts/overrides_kuberblue/usr/lib/systemd/system/crio.service.d/10-exec-path.conf
# Fix for cri-o1.35 RPM: unit file ships with /usr/local/bin/crio but
# binary is installed to /usr/bin/crio by the RPM.
[Service]
ExecStart=
ExecStart=/usr/bin/crio $CRIO_CONFIG_OPTIONS $CRIO_RUNTIME_OPTIONS $CRIO_STORAGE_OPTIONS $CRIO_NETWORK_OPTIONS $CRIO_METRICS_OPTIONS
